# templates/borgmatic.timer.j2
[Unit]
Description=Run configured borgmatic backup

[Timer]
OnCalendar={{ borgmatic_schedule_linux }}
Persistent=true

[Install]
WantedBy=timers.target
